# r.conf
[program:r]
command=python3 /home/curator2/r/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/r/test.err.log
stdout_logfile=/home/curator2/r/test.out.log